# deploy/chromium-reload.service
# Manual Chromium reload — restarts ONLY the browser unit (not labwc, not the
# backend). systemd USER service. Installs to ~/.config/systemd/user/.
#
# No timer: nightly browser hygiene is owned by the power-off window (the smart
# plug cold-boots the box at 06:00 daily, so every day starts with a fresh
# Chromium). This unit remains as the documented way to pick up a same-day
# frontend deploy without waiting for that boot:
#
#   systemctl --user start chromium-reload.service
#
# With the no-cache static headers, the restarted browser fetches the new
# bundle immediately.

[Unit]
Description=Reload the pi-dashboard Chromium kiosk

[Service]
Type=oneshot
ExecStart=systemctl --user restart chromium-kiosk.service
