# deploy/chromium-reload.timer
# Fires chromium-reload.service nightly at 04:00 local. systemd USER timer.
# Installs to ~/.config/systemd/user/.
#
#   systemctl --user enable --now chromium-reload.timer
#
# 04:00 is inside the 1a-6a CSS blackout (reload flash hidden behind the overlay)
# and after the 03:00 unattended-upgrades reboot window, so on an upgrade night
# the box has already rebooted and this is a harmless no-op restart. Persistent=false:
# a reload is pure hygiene, so don't replay a missed run on the next boot (which
# would bounce the browser moments after startup).

[Unit]
Description=Nightly pi-dashboard Chromium reload

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=false

[Install]
WantedBy=timers.target
